# files/system/etc/containers/systemd/wazuh-agent.container
[Unit]
Description=Wazuh Agent
After=network-online.target wazuh-agent-init.service
Wants=network-online.target wazuh-agent-init.service

[Container]
Image=docker.io/wazuh/wazuh-agent:4.13.1
ContainerName=wazuh-agent

# Mount configuration, logs, and queue for persistence
Volume=/var/ossec/etc:/var/ossec/etc:Z
Volume=/var/ossec/logs:/var/ossec/logs:Z
Volume=/var/ossec/queue:/var/ossec/queue:Z

# Network - host mode for agent communication
Network=host

# Run as root (required for security monitoring)
User=root

# Security options
SecurityLabelDisable=false

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target
