# files/system/etc/systemd/system/wazuh-agent-init.service
[Unit]
Description=Initialize Wazuh Agent Configuration
Before=wazuh-agent.service
ConditionPathIsDirectory=/var/ossec/etc
ConditionDirectoryNotEmpty=!/var/ossec/etc

[Service]
Type=oneshot
RemainAfterExit=yes
# Copy default config files from container if /var/ossec/etc is empty
ExecStart=/usr/bin/podman run --rm \
  -v /var/ossec/etc:/target:Z \
  docker.io/wazuh/wazuh-agent:4.13.1 \
  sh -c "cp -a /var/ossec/etc/* /target/"

[Install]
WantedBy=multi-user.target
